# crontab
0 */6 * * * /usr/bin/curl http://127.0.0.1/api/sync/artifacts
0 */12 * * * /usr/bin/curl http://127.0.0.1/api/utils/refresh
